# container/var/lib/dpkg/info/lpr.conffiles
/etc/cron.weekly/lpr
/etc/default/lpd
/etc/init.d/lpd
/etc/logcheck/ignore.d.server/lpr
